# crontab
# Add this line to your crontab, and add the path to this file to the PATH in your crontab

10 8 * * * root backup-volumes.sh
